# config/launchd/com.anima.weight_precache_monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anima.weight_precache_monitor</string>

    <!--
        pre-H100 ROI P4 — weight precache progress + ETA refresh

        Runs tool/weight_precache_progress_monitor.hexa every 600 sec (10 min).
        Reads state/h100_weight_precache_progress.jsonl + manifest, computes
        per-path % done and overall ETA, writes
        state/weight_precache_eta.json (atomic SSOT replacement).

        Pure read+derive+write — no network, no mutation of progress.jsonl,
        no Slack notify. Safe to run while h100_weight_precache.bash is
        actively appending to the progress file.

        ACTIVATION (USER-EXPLICIT ONLY)
            launchctl load -w /Users/ghost/core/anima/config/launchd/com.anima.weight_precache_monitor.plist
        The agent NEVER auto-loads this plist.

        StartInterval = 600 sec  →  10 min
        RunAtLoad     = false
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/.hx/bin/hexa</string>
        <string>run</string>
        <string>/Users/ghost/core/anima/tool/weight_precache_progress_monitor.hexa</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/core/anima</string>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/core/anima/logs/weight_precache_monitor.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/core/anima/logs/weight_precache_monitor.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ghost/.hx/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
