# config/launchd/com.anima.dist_native_build_periodic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anima.dist_native_build_periodic</string>

    <!--
        ROI item T2 — dist/ native build daily cron.

        Runs tool/dist_native_build_periodic.hexa once every 24h. The wrapper
        invokes tool/build_fixup_native.bash --selftest, which rebuilds
        dist/serve_alm_persona.native + dist/h100_post_launch_ingest.native
        against the latest tool/*.hexa sources and verifies each binary's
        own --selftest exit code == 0.

        StartInterval = 86400 sec  →  24 hours.
        RunAtLoad = false  →  first run waits for full interval (avoids
                              boot-time clang contention with H100 launches).

        Append-only audit log: state/dist_native_build_periodic_log.jsonl

        Activation: tool/launchd_install_local.bash  (USER must run manually).
        AGENT NEVER auto-loads this plist.
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/hexa-lang/hexa</string>
        <string>run</string>
        <string>/Users/ghost/core/anima/tool/dist_native_build_periodic.hexa</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/core/anima</string>

    <key>StartInterval</key>
    <integer>86400</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/core/anima/logs/dist_native_build_periodic.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/core/anima/logs/dist_native_build_periodic.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
