# config/launchd/com.anima.cert_dag_periodic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anima.cert_dag_periodic</string>

    <!--
        ROI item E26 — cert_dag orphan check periodic (every 12h).

        Re-runs tool/cert_dag_generator.hexa to refresh state/cert_dag.json
        and surface any orphaned cert artifacts (unreferenced nodes / dangling
        edges) BEFORE they trip the cert_loaded threshold gate during a live
        H100 ingest.

        StartInterval = 43200 sec  →  12 hours.
        RunAtLoad = false.

        Activation: tool/launchd_install_local.bash  (USER must run manually).
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/hexa-lang/hexa</string>
        <string>run</string>
        <string>/Users/ghost/core/anima/tool/cert_dag_generator.hexa</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/core/anima</string>

    <key>StartInterval</key>
    <integer>43200</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/core/anima/logs/cert_dag_periodic.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/core/anima/logs/cert_dag_periodic.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
