# config/launchd/com.anima.h100_auto_kill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anima.h100_auto_kill</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/hexa-lang/hexa</string>
        <string>run</string>
        <string>/Users/ghost/core/anima/tool/h100_auto_kill.hexa</string>
        <string>--dry-run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/core/anima</string>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/core/anima/logs/h100_auto_kill.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/core/anima/logs/h100_auto_kill.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
