# config/launchd/com.anima.auto_evolution.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anima.auto_evolution</string>

    <!--
        P9 — anima auto-evolution 12h cycle
        Spec: docs/anima_proposal_stack_paradigm_20260422.md §3 12h cycle.

        Runs tool/auto_evolution_loop.hexa --apply every 43200s (12 hours),
        which orchestrates the 8-step proposal stack cycle:
          1. proposal_inventory_init
          2. per-pending refinement (version bump)
          3. proposal_emit --apply  (dedup, max 20)
          4. proposal_cluster_detect --apply
          5. proposal_conflict_detect --apply  (auto debate folder routing)
          6. meta/metrics.json refresh
          7. proposal_dashboard regenerate
          8. notify (slack via env SLACK_WEBHOOK_URL when set)

        Per-step results stream to:
          state/proposals/meta/cycle_log.jsonl

        ACTIVATION (USER-EXPLICIT ONLY)
          bash tool/launchd_install_local.bash --activate
        The agent NEVER auto-activates this plist. The user must opt in.

        StartInterval = 43200 sec  →  12 hours.
        RunAtLoad     = false.
    -->

    <key>ProgramArguments</key>
    <array>
        <string>bash</string>
        <string>-c</string>
        <string>/Users/ghost/.hx/bin/hexa run tool/auto_evolution_loop.hexa --apply</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/core/anima</string>

    <key>StartInterval</key>
    <integer>43200</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/anima_auto_evolution.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/anima_auto_evolution.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ghost/.hx/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
