# config/launchd/com.anima.log_rotation_weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.anima.log_rotation_weekly.plist
  ROI V4 (pre-H100): weekly zstd rotation of /tmp/*.log + logs/*.log + state/*.log
  Schedule: every Sunday at 04:15 local time (matches storage_roi_audit.json
  G40 cron_recommendation "15 4 * * 0").

  ACTIVATION (manual; tool does NOT auto-load this plist):
    cp config/launchd/com.anima.log_rotation_weekly.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.anima.log_rotation_weekly.plist
    launchctl enable      gui/$(id -u)/com.anima.log_rotation_weekly
    launchctl kickstart -k gui/$(id -u)/com.anima.log_rotation_weekly   # force one immediate run for verification

  STATUS / VERIFY:
    launchctl print gui/$(id -u)/com.anima.log_rotation_weekly | head -40
    cat /Users/ghost/core/anima/state/log_rotation_zstd_last.json

  DEACTIVATION:
    launchctl bootout gui/$(id -u)/com.anima.log_rotation_weekly
    rm ~/Library/LaunchAgents/com.anima.log_rotation_weekly.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anima.log_rotation_weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ghost/core/anima/tool/log_rotation_zstd.bash</string>
        <string>--apply</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/core/anima</string>

    <!-- Sundays 04:15 local; matches G40 cron_recommendation -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/core/anima/logs/log_rotation_weekly.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/core/anima/logs/log_rotation_weekly.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>ANIMA_ROOT</key>
        <string>/Users/ghost/core/anima</string>
    </dict>

    <!-- Resource bounds: zstd -19 -T0 is CPU-heavy; nice it down. -->
    <key>Nice</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
