# config/launchd/com.anima.airgenome_keyword_dispatch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anima.airgenome_keyword_dispatch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/core/anima/bin/airgenome_keyword_dispatch</string>
        <string>--dry-run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/core/anima</string>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/core/anima/logs/airgenome_keyword_dispatch.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/core/anima/logs/airgenome_keyword_dispatch.err.log</string>
</dict>
</plist>
