# config/launchd/com.anima.lint_cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anima.lint_cron</string>

    <!--
        Daily anima F1-F6 lint sweep (raw#7 + raw#8 residual tracker).

        Runs the cross-repo hexa-lang ai_native_scan against the anima
        tree and appends a single JSON line to state/lint_cron_history.jsonl
        capturing the per-rule counts. The exemption files
        (.raw-exemptions/raw_6.list, raw_7.list, raw_8.list) declare which
        violations are jurisdiction-out (research/bench/experimental/deploy);
        the cron sweep tracks the residual delta over time so re-justification
        at expiry (2027-04-26) has a hard data trail.

        StartCalendarInterval = daily 03:17 KST (UTC 18:17 prev day).
        RunAtLoad = false.

        SCHEDULE RATIONALE (raw#83 — schedule the worst-of-day for cron):
          03:17 = overnight quiet window, ~14h after typical commit burst,
          long enough for any in-flight refactor to land and short enough
          that drift surfaces within a working day.

        Activation (manual — no auto-bootstrap):
          1) cp config/launchd/com.anima.lint_cron.plist ~/Library/LaunchAgents/
          2) launchctl bootstrap gui/$UID ~/Library/LaunchAgents/com.anima.lint_cron.plist
          3) launchctl enable gui/$UID/com.anima.lint_cron
          4) (optional kick) launchctl kickstart -k gui/$UID/com.anima.lint_cron

        Tear-down:
          launchctl bootout gui/$UID/com.anima.lint_cron
          rm ~/Library/LaunchAgents/com.anima.lint_cron.plist

        Output:
          state/lint_cron_history.jsonl  — one JSON per run
            { ts, total, f1, f2, f3, f4, f5, f6, exit_code }
          logs/lint_cron.out.log  — stdout (raw scanner output)
          logs/lint_cron.err.log  — stderr (scan errors)
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/ghost/core/anima &amp;&amp; OUT="$(/Users/ghost/core/hexa-lang/hexa /Users/ghost/core/hexa-lang/tool/ai_native_scan.hexa 2&gt;&amp;1)"; EC=$?; TS="$(date -u +%Y-%m-%dT%H:%M:%SZ)"; F1=$(printf '%s' "$OUT" | grep -c '^  F1-' || true); F2=$(printf '%s' "$OUT" | grep -c '^  F2-' || true); F3=$(printf '%s' "$OUT" | grep -c '^  F3-' || true); F4=$(printf '%s' "$OUT" | grep -c '^  F4-' || true); F5=$(printf '%s' "$OUT" | grep -c '^  F5-' || true); F6=$(printf '%s' "$OUT" | grep -c '^  F6-' || true); TOTAL=$((F1+F2+F3+F4+F5+F6)); printf '{"ts":"%s","total":%d,"f1":%d,"f2":%d,"f3":%d,"f4":%d,"f5":%d,"f6":%d,"exit_code":%d}\n' "$TS" "$TOTAL" "$F1" "$F2" "$F3" "$F4" "$F5" "$F6" "$EC" &gt;&gt; /Users/ghost/core/anima/state/lint_cron_history.jsonl; printf '%s\n' "$OUT"</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/core/anima</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>17</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/core/anima/logs/lint_cron.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/core/anima/logs/lint_cron.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HEXA_LOCAL</key>
        <string>1</string>
    </dict>
</dict>
</plist>
