# config/launchd/com.anima.runpod_credit_check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anima.runpod_credit_check</string>

    <!--
        pre-H100 ROI P3 — runpod credit balance guard

        Runs tool/runpod_credit_check.hexa every 300 sec (5 min).
        Default invocation is dry — probes balance, writes
        state/runpod_credit_status.json. If balance < $1000 the SSOT carries
        alert=true, but Slack notify only fires when this plist's
        ProgramArguments includes --apply AND env SLACK_WEBHOOK_URL is set.

        ACTIVATION (USER-EXPLICIT ONLY)
            launchctl load -w /Users/ghost/core/anima/config/launchd/com.anima.runpod_credit_check.plist
        The agent NEVER auto-loads this plist.

        StartInterval = 300 sec  →  5 min
        RunAtLoad     = false
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/.hx/bin/hexa</string>
        <string>run</string>
        <string>/Users/ghost/core/anima/tool/runpod_credit_check.hexa</string>
        <string>--apply</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/core/anima</string>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/core/anima/logs/runpod_credit_check.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/core/anima/logs/runpod_credit_check.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ghost/.hx/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
