# config/launchd/com.anima.adversarial_bench_periodic.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anima.adversarial_bench_periodic</string>

    <!--
        ROI item D24 — adversarial_bench periodic cron (every 6h).

        Re-runs tool/adversarial_bench.hexa to produce a fresh
        state/adversarial_bench_last.json and detect any drift in raw#12
        cherry-pick immunity (3/3 flip differentiation).

        StartInterval = 21600 sec  →  6 hours.
        RunAtLoad = false  →  first run waits for full interval (avoids
                              boot-time CPU contention with H100 launches).

        Activation: tool/launchd_install_local.bash  (USER must run manually).
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/Dev/hexa-lang/hexa</string>
        <string>run</string>
        <string>/Users/ghost/core/anima/tool/adversarial_bench.hexa</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/core/anima</string>

    <key>StartInterval</key>
    <integer>21600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/ghost/core/anima/logs/adversarial_bench_periodic.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ghost/core/anima/logs/adversarial_bench_periodic.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
