# anima-os/launchd/roadmap_live_daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anima.roadmap-live-daemon</string>

    <!--
      anima roadmap_live_daemon — OS-level .roadmap live enforcement.
      Replaces git pre-commit hook (폐기됨): hooks bypass is OK,
      launchd periodic tick is the SSOT gate (raw#1/10/15).

      Each tick:
        1. chflags uchg re-apply on .roadmap (raw#1)
        2. roadmap_lint (hexa-lang canonical tool)
        3. β coverage audit (file exists + parse)
        4. on drift: auto_reflect.hexa --dry-run (report-only)
        4b. roadmap_integrity_guard (7-check: lint/uchg/paths/file/ssot/lock/coverage)
            → shared/state/integrity_drift_<ts>.json on WARN/FAIL
            → shared/state/roadmap_integrity_result.json (rolling SSOT)
        4c. phase_progression_controller --stage live  (3-stage OS pipeline
            Stage 3 hook). Runs phase_live_adjuster for WARN/drift/debt
            auto-fix within spec bounds per live_adjust_policy.
            → shared/state/phase_live_adjust_log.json (JSONL)
            → shared/state/phase_progression_cert.json (rolling cert)
        5. .raw-audit hash-chain append (raw#10)

      Triggers (constraint #1 resolution — fsevents stdlib binding 不要):
        - StartInterval = 3600 (1h periodic backstop, RunAtLoad = true)
        - WatchPaths    = .roadmap / edu/paths.json / .ext  (real-time)
        - ThrottleInterval = 10s (burst suppression: rapid edits → 1 tick)

      WatchPaths wakes the daemon the moment launchd's FSEvent subsystem
      observes a change under the listed paths. Native fsevents bindings
      in hexa are therefore unnecessary for the MVP; launchd handles
      kernel→userspace wake-up via its own subsystem.
    -->

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ghost/.hx/bin/hexa</string>
        <string>/Users/ghost/core/anima/tool/roadmap_live_daemon.hexa</string>
        <string>--tick</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ghost/core/anima</string>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>WatchPaths</key>
    <array>
        <string>/Users/ghost/core/anima/.roadmap</string>
        <string>/Users/ghost/core/anima/edu/paths.json</string>
        <string>/Users/ghost/core/anima/.ext</string>
    </array>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/roadmap_live_daemon.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/roadmap_live_daemon.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/ghost/.hx/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/ghost</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>15</integer>

    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
